# systemd/luanti-backup.service
[Unit]
Description=Luanti daily backup

[Service]
Type=oneshot
ExecStart=/home/user/scripts/luanti_backup.sh