# systemd/luanti-backup.timer
[Unit]
Description=Daily Luanti backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target